# deploy/api.service
[Unit]
Description=A股股价看板 API（FastAPI + DuckDB，同进程托管前端 dist）
After=network.target redis-server.service
Wants=redis-server.service
StartLimitIntervalSec=0

[Service]
Type=simple
User=root
WorkingDirectory=/root/Quantitative_A-shares/WechatNum
Environment="DUCKDB_PATH=/root/Quantitative_A-shares/WechatNum/market.duckdb"
Environment="REDIS_URL=redis://127.0.0.1:6379/0"
Environment="DUCKDB_MEMORY_LIMIT=400MB"
ExecStart=/root/.local/bin/uv run uvicorn src.api.main:app --host 0.0.0.0 --port 8501 --workers 1
Restart=on-failure
RestartSec=10s

[Install]
WantedBy=multi-user.target
